# AGENT/agent.mon/bin/falcon_migrate_ecs_tag_v2.crontab
*/5 * * * * root flock -n /tmp/falcon_migrate_ecs_tag_v2.lock /opt/mydan/dan/agent.mon/bin/falcon_migrate_ecs_tag_v2 > /tmp/falcon_migrate_ecs_tag_v2.log 2>&1
