# AGENT/monitor/whiteking.agent.mointor.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
* * * * * root flock -w 1 -x /tmp/whiteking.agent.mointor.lock /data/Software/mydan/AGENT/monitor/falcon_push.sh >/dev/null 2>&1
