# Connector/openc3.crontab
CRON_TZ=Asia/Shanghai
TZ=Asia/Shanghai

15 2 * * * root /usr/bin/c3mc-ci-gitreport-cron  > /tmp/c3mc-ci-gitreport-cron.log  2>&1
45 3 * * * root /usr/bin/c3mc-mon-monreport-cron > /tmp/c3mc-mon-monreport-cron.log 2>&1
1  3 * * * root /usr/bin/c3mc-k8s-backup-cron    > /tmp/c3mc-k8s-backup-cron.log    2>&1
5  5 * * * root /usr/bin/c3mc-sys-backup         > /tmp/c3mc-sys-backup.log         2>&1
15 4 * * * root /usr/bin/c3mc-mon-nlow-cron      > /tmp/c3mc-mon-nlow-cron.log      2>&1
1  * * * * root /usr/bin/c3mc-device-timemachine > /tmp/c3mc-device-timemachine.log      2>&1

*/3 * * * * root /usr/bin/c3mc-device-accountdb-dump > /tmp/c3mc-device-accountdb-dump.log 2>&1

*/5 * * * * root flock -n /tmp/c3mc-cloud-update-cloud-account.lock /usr/bin/c3mc-cloud-update-cloud-account > /tmp/c3mc-cloud-update-cloud-account.log 2>&1

*/2 * * * * root /usr/bin/c3mc-device-bind-info-cache set > /tmp/c3mc-device-bind-info-cache.log 2>&1
*/5 * * * * root /usr/bin/c3mc-device-base-info-cache set > /tmp/c3mc-device-base-info-cache.log 2>&1
*/5 * * * * root /usr/bin/c3mc-device-cat-all-cache   set > /tmp/c3mc-device-cat-all-cache.log 2>&1
*/5 * * * * root /usr/bin/c3mc-device-nodeinfo-allnode-cache set > /tmp/c3mc-device-nodeinfo-allnode-cache.log 2>&1
1 * * * * root /data/Software/mydan/AGENT/device/bin/depttree > /tmp/depttree.log 2>&1

* * * * * root flock -n /tmp/c3mc-mon-proxy-metricsd.lock /usr/bin/c3mc-mon-proxy-metricsd > /tmp/c3mc-mon-proxy-metricsd.log 2>&1

*/5 * * * * root /usr/bin/c3mc-mon-make-uexip > /data/Software/mydan/AGENT/conf/openc3_uexip.yml.temp.$$ && mv /data/Software/mydan/AGENT/conf/openc3_uexip.yml.temp.$$ /data/Software/mydan/AGENT/conf/openc3_uexip.yml

*/5 * * * * root flock -n /tmp/c3mc-device-account-proxyd.lock  /usr/bin/c3mc-device-account-proxyd > /tmp/c3mc-device-account-proxyd.log 2>&1

*/5 * * * * root /data/Software/mydan/Connector/pp/crontask/cron-make > /tmp/cron-make.log 2>&1

1 * * * * root /data/Software/mydan/Connector/pp/mmon/resourcelow/makecache >/tmp/resourcelow-makecache.loc 2>&1
* * * * * root flock -n /tmp/serialcall.lock /data/Software/mydan/Connector/pp/mmon/serialcall/run.sh > /tmp/serialcall.log 2>&1

*/30 * * * * root /data/Software/mydan/Connector/pp/device/buildintree/run > /tmp/c3.buildintree.log 2>&1

*/10 * * * * root flock -n /tmp/alertmanagercrontab.lock /data/Software/mydan/Connector/pp/mmon/alertmanagercrontab/run.sh > /tmp/alertmanagercrontab.log 2>&1

0 10 * * * root /data/Software/mydan/Connector/pp/device/crontab/ssl-certificate-expired-notify/run.sh

*/3 * * * * root flock -n /tmp/cislave.sync1.lock /data/Software/mydan/CI/cislave/sync/sync1.sh  > /tmp/cislave.sync1.log 2>&1
*/5 * * * * root flock -n /tmp/cislave.sync1.lock /data/Software/mydan/CI/cislave/sync/sync2.sh  > /tmp/cislave.sync2.log 2>&1

* * * * * root flock -n /tmp/dbtask.lock /data/Software/mydan/Connector/pp/dbfix/dbtask.sh  > /tmp/dbtask.log 2>&1

*/10 * * * * root flock -n /tmp/mmon-node-alive-make-exip.lock /data/Software/mydan/Connector/pp/mmon/node-alive/make-exip > /tmp/mmon-node-alive-make-exip.log 2>&1
* * * * * root flock -n /tmp/mmon-node-alive.lock /data/Software/mydan/Connector/pp/mmon/node-alive/run.sh > /tmp/mmon-node-alive.log 2>&1

*/5 * * * * root flock -n /tmp/resources-fee-change-notify.lock /data/Software/mydan/Connector/pp/device/crontab/resources-fee-change-notify/notifyd  > /tmp/resources-fee-change-notify.log 2>&1

* * * * * root flock -n /tmp/c3mc-mon-exmesg-cron.lock  /usr/bin/c3mc-mon-exmesg-cron   > /tmp/c3mc-mon-exmesg-cron.log      2>&1
* * * * * root flock -n /tmp/c3mc-mon-exmesg-maker-cron.lock  /usr/bin/c3mc-mon-exmesg-maker-cron   > /tmp/c3mc-mon-exmesg-maker-cron.log      2>&1

1 */2 * * * root flock -n /tmp/service-analysis-make.lock /data/Software/mydan/Connector/pp/service-analysis/make.sh > /tmp/service-analysis-make.log 2>&1
15 */2 * * * root flock -n /tmp/service-analysis-nginx-make.lock /data/Software/mydan/Connector/pp/service-analysis/nginx/sync > /tmp/service-analysis-nginx-make.log 2>&1

* * * * * root flock -n /tmp/http-monitor.lock /data/Software/mydan/Connector/pp/mmon/http-monitor/run > /tmp/http-monitor.log 2>&1

*/5 * * * * root flock -n /tmp/find-v3-data-make.lock /data/Software/mydan/Connector/pp/device/find-v3-data-make > /tmp/find-v3-data-make.log 2>&1
