# Connector/openc3.crontab
15 2 * * * root /data/Software/mydan/CI/server/gitreport 2>&1 > /tmp/gitreport.log |/data/Software/mydan/CI/pp/c3_ci_pp_addtime >> /var/log/open-c3.error.log
30 * * * * root /data/Software/mydan/CI/server/flowreport > /tmp/flowreport.log 2>&1
