# Connector/openc3.crontab
15 2 * * * root /data/Software/mydan/CI/server/gitreport > /tmp/gitreport.log 2>&1
