# Connector/openc3.crontab
15 2 * * * root /data/Software/mydan/CI/server/gitreport > /tmp/gitreport.log 2>&1
30 * * * * root /data/Software/mydan/CI/server/flowreport > /tmp/flowreport.log 2>&1
