# Connector/openc3.crontab
15 2 * * * root /usr/bin/c3mc-ci-gitreport-cron > /tmp/c3mc-ci-gitreport-cron.log 2>&1
1 3 * * * root /usr/bin/c3mc-k8s-backup-cron > /tmp/c3mc-k8s-backup-cron.log 2>&1
