# Connector/spotx.crontab
* * * * * root flock -xn /var/lock/spotx.lock /data/Software/mydan/Connector/spotx.sh > /tmp/spotx.log 2>&1
